# .config/systemd/user/pipewire-filters.target
[Unit]
Description=Pipewire plugins
PartOf=wireplumber.service
After=wireplumber.service

[Install]
WantedBy=default.target
