# .config/systemd/user/pipewire-filters.target
[Unit]
Description=Pipewire plugins

[Install]
WantedBy=default.target
