# .config/systemd/user/wireplumber.service.d/override.conf
[Service]
Environment=WIREPLUMBER_DEBUG=audio-links:3
